# scripts/omniroute.service
[Unit]
Description=OmniRoute AI Router
After=network.target

[Service]
Type=simple
User=openclaw
WorkingDirectory=/home/openclaw/.openclaw/workspace
ExecStart=/home/openclaw/.npm-global/bin/omniroute --no-open --port 20128
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=HOME=/home/openclaw
Environment=PATH=/home/openclaw/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
# IMPORTANT: Replace YOUR_ANTIGRAVITY_OAUTH_CLIENT_SECRET_HERE with your actual secret
# Get it from: Google Cloud Console → APIs & Services → Credentials → OAuth 2.0 Client IDs
Environment=ANTIGRAVITY_OAUTH_CLIENT_SECRET=YOUR_ANTIGRAVITY_OAUTH_CLIENT_SECRET_HERE

[Install]
WantedBy=multi-user.target
